# docker-image/kubectl/skelton/etc/profile.d/kubectl.sh
#!/bin/bash


#
# kubectl settings
#
# setup autocomplete in bash, bash-completion package should be installed first.
source <(kubectl completion bash)
